# nsite/dial-a-geddon.upstart.conf
#
#  dial-a-geddon.conf -- ubuntu upstart file to run dial-a-geddon nodejs site
#
# This is an Ubuntu Upstart .conf file, for boot-up starting of node server,
# not some other kind of .conf file for dag.
# must have a copy in /etc/init that points to this
#   sudo cp -i dial-a-geddon.upstart.conf /etc/init

description "dial-a-geddon node.js server"
# adapted from author      "kvz - http://kevin.vanzonneveld.net"
author      "Allan Bonadio"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/home/www"

	# run as www
    exec sudo -u www /usr/local/bin/node /dvl/dial-a-geddon/nsite/app.js >> /var/log/dial-a-geddon.log 2>&1
end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
   echo `date +%c` "starting nodejs for dag" >> /var/log/dial-a-geddon.log
end script
